# src/lib/systemd/system/minecraft-scheduled-restart@.timer
[Unit]
Description=Minecraft Server Scheduled Restart timer (CyanMC) %i

[Timer]
OnCalendar=*-*-* 00:00:00 UTC

[Install]
WantedBy=multi-user.target
