# systemd/catstar-backup.service
[Unit]
Description=Catstar backup service

[Service]
WorkingDirectory=/root
ExecStart=/bin/bash backup.sh
Environment=HOME=/root
Type=oneshot
